# rootfs/etc/s6-overlay/s6-rc.d/openchamber/finish
#!/bin/bash
# Run cleanup regardless of exit reason
s6-setuidgid agent openchamber stop 2>/dev/null || true
# $1 = exit code; 256 = intentional takedown, not a crash
[ "${1}" -eq 256 ] && exit 0
echo "openchamber crashed (exit ${1}), halting container..."
exec s6-svscanctl -t /run/s6/service
