# rootfs/etc/s6-overlay/s6-rc.d/openchamber/finish
#!/command/execlineb -P
with-contenv
cd /home/agent
openchamber stop
