# rootfs/etc/s6-overlay/s6-rc.d/dockerd-ready/up
#!/command/execlineb -P
with-contenv
/bin/bash -c "for i in $(seq 1 30); do [ -S /var/run/docker.sock ] && exit 0; sleep 1; done; exit 1"
