# rootfs/etc/s6-overlay/s6-rc.d/dockerd/finish
#!/bin/bash
[ "${1}" -eq 256 ] && exit 0
echo "dockerd crashed (exit ${1}), halting container..."
exec /command/s6-svscanctl -t /run/s6/service
