# rootfs/etc/s6-overlay/s6-rc.d/dockerd/finish
#!/bin/bash
exec /run/s6/basedir/bin/halt
