# rootfs/etc/s6-overlay/s6-rc.d/dockerd/run
#!/command/execlineb -P
fdmove -c 2 1
with-contenv
/usr/bin/dockerd
