# rootfs/etc/s6-overlay/s6-rc.d/opencode/finish
#!/bin/bash
[ "${1}" -eq 256 ] && exit 0
echo "opencode crashed (exit ${1}), halting container..."
exec s6-svscanctl -t /run/s6/service
