# rootfs/etc/s6-overlay/s6-rc.d/agent-init/up
#!/command/execlineb -P
with-contenv
export HOME /home/agent
s6-setuidgid agent
cd /home/agent
/usr/local/bin/agent-setup.sh
