# rootfs/etc/s6-overlay/s6-rc.d/agent-init/up
#!/command/execlineb -P
with-contenv
if { /bin/sh -c "/usr/local/bin/generate-java-cacerts.sh" }
export HOME /home/agent
s6-setuidgid agent
cd /home/agent
/usr/local/bin/agent-setup.sh
