# rootfs/etc/s6-overlay/s6-rc.d/agent-init/up
#!/command/execlineb -P
with-contenv
cd /home/agent
/usr/local/bin/agent-setup.sh
